# layout/DEBIAN/postinst
#!/bin/sh

launchctl stop com.zx02.xenlived
launchctl unload /Library/LaunchDaemons/com.zx02.xenlived.plist
launchctl load /Library/LaunchDaemons/com.zx02.xenlived.plist